# sparse/etc/systemd/user/ngfd.service.d/wait.conf
[Service]
ExecStartPre=/bin/sh -c "while [ ! -f /run/user-session/droid-misc.flag ]; do sleep 1; done"

